# systemd/home-worker-ota-deps-build@.service
[Unit]
Description=Build Home Worker OTA native dependencies for operation %i
ConditionPathExists=/run/home-worker/prepare/%i.json
ConditionPathExists=/run/home-worker/prepare/%i/candidate
PartOf=home-worker-ota-prepare@%i.service
BindsTo=home-worker-ota-prepare@%i.service

[Service]
Type=oneshot
User=homeworker
Group=homeworker
UMask=0077
WorkingDirectory=/run/home-worker/prepare/%i/candidate
ExecStart=/usr/bin/node /usr/lib/home-worker/ota-prepare.mjs build %i
KillMode=control-group
TimeoutStartSec=20min
TimeoutStopSec=30s
MemoryMax=300M
TasksMax=64
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RemoveIPC=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
InaccessiblePaths=/opt/home-worker/shared
InaccessiblePaths=/opt/home-worker/data
InaccessiblePaths=-/opt/home-worker/.env
InaccessiblePaths=-/opt/home-worker/.npmrc
InaccessiblePaths=-/opt/home-worker/.yarnrc.yml
InaccessiblePaths=-/opt/home-worker/releases/.npmrc
InaccessiblePaths=-/opt/home-worker/releases/.yarnrc.yml
InaccessiblePaths=/home/homeworker/.pm2
InaccessiblePaths=-/etc/npmrc
InaccessiblePaths=-/etc/yarnrc.yml
ReadWritePaths=/run/home-worker/prepare/%i/candidate /run/home-worker/prepare/%i/tmp
Environment=PATH=/usr/bin:/bin
UnsetEnvironment=HOME XDG_CONFIG_HOME XDG_CACHE_HOME NODE_OPTIONS npm_config_userconfig npm_config_registry npm_config_proxy npm_config_https_proxy NPM_CONFIG_USERCONFIG NPM_CONFIG_REGISTRY NPM_CONFIG_PROXY NPM_CONFIG_HTTPS_PROXY YARN_NPM_AUTH_TOKEN YARN_NPM_AUTH_IDENT YARN_HTTP_PROXY YARN_HTTPS_PROXY HTTP_PROXY HTTPS_PROXY ALL_PROXY NO_PROXY
StandardOutput=journal
StandardError=journal

[Install]
# Coordinator-started only; this template is never enabled or sudo-startable.
