# systemd/home-worker-ota-prepare@.service
[Unit]
Description=Prepare Home Worker OTA dependencies for operation %i
ConditionPathExists=/run/home-worker/prepare/%i.json
ConditionPathExists=/run/home-worker/prepare/%i/candidate

[Service]
Type=oneshot
User=homeworker
Group=homeworker
UMask=0077
WorkingDirectory=/run/home-worker/prepare/%i/candidate
ExecStart=/usr/bin/flock --no-fork --exclusive /run/home-worker/ota-prepare.lock /usr/bin/node /usr/lib/home-worker/ota-prepare.mjs %i
KillMode=control-group
TimeoutStartSec=30min
TimeoutStopSec=30s
MemoryMax=512M
TasksMax=64
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RemoveIPC=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
InaccessiblePaths=/opt/home-worker/shared
InaccessiblePaths=/opt/home-worker/data
InaccessiblePaths=-/opt/home-worker/.env
InaccessiblePaths=/home/homeworker/.pm2
ReadWritePaths=/run/home-worker/prepare/%i/candidate /run/home-worker/prepare/%i/tmp
Environment=HOME=/run/home-worker/prepare/%i/tmp/home
Environment=TMPDIR=/run/home-worker/prepare/%i/tmp
Environment=TMP=/run/home-worker/prepare/%i/tmp
Environment=TEMP=/run/home-worker/prepare/%i/tmp
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512
Environment=npm_config_jobs=1
Environment=JOBS=1
Environment=YARN_ENABLE_NETWORK=false
Environment=YARN_ENABLE_GLOBAL_CACHE=false
Environment=YARN_ENABLE_IMMUTABLE_INSTALLS=true
Environment=YARN_ENABLE_IMMUTABLE_CACHE=true
Environment=YARN_IGNORE_PATH=true
StandardOutput=journal
StandardError=journal

[Install]
# Instances are started on demand only; this template is never enabled.
